# supervisor/schedule.conf
[program:scheduler]
directory=/var/www/html
command=/usr/local/bin/php /var/www/html/artisan schedule:work
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
;stdout_logfile=/dev/stdout
'stdout_logfile_maxbytes=0
;stderr_logfile=/dev/stderr
;stderr_logfile_maxbytes=0
stderr_logfile=/var/log/supervisor/scheduler.err.log
stdout_logfile=/var/log/supervisor/scheduler.out.log
stopsignal=QUIT
